# crontab
# Run script.rb at midnight every day
0 0 * * * /app/run_script.sh >> /var/log/cron.log 2>&1
# Ensure the cron daemon does not exit
# An empty line is required at the end of this file for a valid cron file.
